# files/etc/uci-defaults/99-custom-ip
#!/bin/sh
uci -q batch << EOF
set network.lan.ipaddr='192.168.2.1'
set network.lan.netmask='255.255.255.0'
commit network
EOF
exit 0
